# systemd/simbox-monthly-sms.timer
[Unit]
Description=Monthly keep-alive SMS timer for sim-box

[Timer]
OnCalendar=*-*-01 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
