# systemd/simbox-reboot.timer
[Unit]
Description=Daily reboot timer for sim-box

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
