# systemd/simbox-monthly-sms.service
[Unit]
Description=Sim-box monthly keep-alive SMS
After=network.target simbox-modem-init.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/simbox-monthly-sms.sh
